# _integrations/qdrant/cron/help.{{PROJECT_SLUG}}.qdrant-weekly-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>help.{{PROJECT_SLUG}}.qdrant-weekly-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{PROJECT_ROOT}}/_integrations/qdrant/cron/run-weekly-sync.sh</string>
    </array>

    <!-- Sunday (Weekday 0) at 22:00 local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- If the machine is asleep at the scheduled time, launchd reruns at wake -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{{PROJECT_ROOT}}/_integrations/qdrant/logs/cron-weekly.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>{{PROJECT_ROOT}}/_integrations/qdrant/logs/cron-weekly.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>{{PROJECT_ROOT}}/_integrations/qdrant</string>
</dict>
</plist>
